# backend/scripts/sync-clients-table.timer
[Unit]
Description=Periodically sync the 3x-ui v3 clients table (desync prevention)
Documentation=https://github.com/rsidorenko/dev-sandbox backend/scripts/sync_clients_table.py

[Timer]
# First run ~2 min after boot, then every 5 minutes. x-ui restart happens only
# when the table changed, so steady-state runs are cheap and non-disruptive.
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
